# tests/packaging/supervisord.conf
[program:tyk-ai-studio]
command=/bin/bash -c 'set -a; source /etc/default/tyk-ai-studio; set +a; exec /opt/tyk-ai-studio/tyk-ai-studio'
user=tyk
directory=/opt/tyk-ai-studio
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/tyk-ai-studio.stdout.log
stderr_logfile=/var/log/tyk-ai-studio.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:tyk-microgateway]
command=/bin/bash -c 'set -a; source /etc/default/tyk-microgateway; set +a; exec /opt/tyk-microgateway/tyk-microgateway'
user=tyk
directory=/opt/tyk-microgateway
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/tyk-microgateway.stdout.log
stderr_logfile=/var/log/tyk-microgateway.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
